# .config/systemd/user/notebook.service
[Unit]
Description=Jupyter

[Service]
ExecStart=jupyter lab  --NotebookApp.notebook_dir=/home/aziz

[Install]
WantedBy=default.target
